# root/etc/services.d/ambari-agent/run
#!/usr/bin/with-contenv sh

ambari-agent start
tail -f /var/log/ambari-agent/ambari-agent.log