# root/etc/services.d/ambari-server/run
#!/usr/bin/with-contenv sh

ambari-server start
tail -f /dev/null